# services/worker/launchd/com.karikris.butterflylens.worker.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.karikris.butterflylens.worker</string>
  <key>Program</key>
  <string>__PYTHON__</string>
  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>-m</string>
    <string>butterflylens_worker.service</string>
    <string>--environment-file</string>
    <string>__ENVIRONMENT_FILE__</string>
    <string>--state-dir</string>
    <string>__STATE_DIR__</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__REPOSITORY__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>__PYTHONPATH__</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>ProcessType</key>
  <string>Background</string>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>ExitTimeOut</key>
  <integer>60</integer>
  <key>Umask</key>
  <string>0077</string>
  <key>StandardOutPath</key>
  <string>__STDOUT_LOG__</string>
  <key>StandardErrorPath</key>
  <string>__STDERR_LOG__</string>
</dict>
</plist>
